# relay/systemd/enclave-api-relay.service
[Unit]
Description=Enclave API relay (fleet discovery + placement)
After=network-online.target
Wants=network-online.target
# never let a crash-loop permanently stop a shared relay (fix 10)
StartLimitIntervalSec=0

[Service]
EnvironmentFile=/etc/nan-relay/api-relay.env
WorkingDirectory=/opt/nan-relay
ExecStart=/usr/bin/node /opt/nan-relay/api-relay.js
Restart=always
RestartSec=2
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
# Writable state for the account/billing modules (accounts.json, orders.json,
# indexer cursor, OFAC cache, ES256 session key). systemd creates
# /var/lib/enclave-relay owned by the dynamic user and exports
# $STATE_DIRECTORY, which store.js picks up; without it (or AUTH_DATA_DIR)
# accounts + billing stay disabled and the relay is a pure router as before.
StateDirectory=enclave-relay
ProtectHome=yes
PrivateTmp=yes
# resource caps (fix 10)
MemoryMax=768M
TasksMax=512
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
